# templates/systemd-restart.path.j2
# Ansible managed. role: {{ role_name }}. template: systemd-restart.path.j2

[Path]
PathModified={{ item.pathmodified }}

[Install]
WantedBy=default.target
